# s6/services.d/docker-gen/run
#!/usr/bin/with-contenv sh

exec docker-gen -watch -notify /usr/bin/reload-haproxy /etc/haproxy/haproxy.cfg.tmpl /etc/haproxy/haproxy.cfg -wait 5 -notify-output